# lib/cron/indefinido.cron
#!/usr/bin/env bash

# crontab:
# * * * * * /bin/sh -l -c 'cd /home/lucca/www/kuva/current && ./lib/cron/pedrocinefoto.cron >> log/search.log 2>> log/search.log'

cd /home/lucca/www/kuva/current

# load rvm ruby
source /usr/local/rvm/environments/ruby-1.9.3-p194

date
echo

bundle exec thor orders:search -e indefinido >> log/search.indefinido.log 2>> log/search.indefinido.log